# debian/accomplishments-daemon.user.service
[Service]
ExecStart=/usr/bin/accomplishments-daemon
